# thetaService.service
[Unit]
Description=StartFuelin

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/fuelup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
